# message_workers/workers_available/worker_broker.conf
[program:worker_broker]
autorestart = true
numprocs = 2
startretries = 3
redirect_stderr = True
stopwaitsecs = 10
process_name = worker_broker_%(process_num)s
priority = 777
directory = /var/lib/databank/message_workers/
command = /var/lib/databank/message_workers/broker.py %(process_num)s
autostart = true
startsecs = 10
stdout_logfile = /var/log/databank/broker.log

